# pv-automation.service
[Unit]
Description=PV-System Automation Daemon (Observer→Engine→Actuator)
After=network-online.target pv-collector.service
Wants=network-online.target
Conflicts=pv-observer.service
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
User=admin
WorkingDirectory=/home/admin/Dokumente/PVAnlage/pv-system
ExecStart=/usr/bin/python3 -m automation.engine.automation_daemon --log-level INFO
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
