# pv-observer.service
[Unit]
Description=PV-System Automation Observer (3-Tier Beobachtungsservice)
After=network-online.target pv-collector.service
Wants=network-online.target
Conflicts=pv-automation.service
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
User=admin
WorkingDirectory=/home/admin/Dokumente/PVAnlage/pv-system
ExecStart=/usr/bin/python3 -m automation.engine.observer --log-level INFO
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
# RAM-DB in tmpfs — kein Problem bei Neustart
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
