# pv-shutdown-persist.service
[Unit]
Description=PV-System: Persist RAM-DB vor Shutdown/Reboot
DefaultDependencies=no
After=local-fs.target
Before=shutdown.target reboot.target halt.target poweroff.target

[Service]
Type=oneshot
User=admin
WorkingDirectory=/srv/pv-system
ExecStart=/bin/bash /srv/pv-system/scripts/persist_on_shutdown.sh
TimeoutStartSec=120

[Install]
WantedBy=halt.target reboot.target poweroff.target
